# ansible/files/supervisord/conf.d/cloudflared.ini
[program:cloudflared]
command={{ user_home }}/cloudflared tunnel run --token {{ vault_cloudflared_token }}
user=dotanarchy
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile={{ user_home }}/.supervisor/logs/cloudflared.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile={{ user_home }}/.supervisor/logs/cloudflared.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
