# ansible/files/supervisord/conf.d/copyparty.ini
[program:copyparty]
command={{ user_home }}/copyparty/virtualenv/bin/python {{ user_home }}/copyparty/copyparty-sfx.py -c {{ user_home }}/copyparty/copyparty.conf
directory={{ user_home }}/copyparty
user=dotanarchy
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile={{ user_home }}/.supervisor/logs/copyparty.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile={{ user_home }}/.supervisor/logs/copyparty.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
