# ansible/files/supervisord/conf.d/deluged.ini
[program:deluged]
command=/usr/bin/deluged
user=dotanarchy
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile={{ user_home }}/.supervisor/logs/deluged.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile={{ user_home }}/.supervisor/logs/deluged.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
