# ansible/files/supervisord/conf.d/deluge-web.ini
[program:deluge-web]
command=/usr/bin/deluge-web
user=dotanarchy
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile={{ user_home }}/.supervisor/logs/deluge-web.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile={{ user_home }}/.supervisor/logs/deluge-web.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
