# ansible/files/supervisord/conf.d/gatus.ini
[program:gatus]
command=bash -c "cd {{ user_home }}/gatus && podman-compose down 2>/dev/null; podman-compose up -d; exec sleep infinity"
directory={{ user_home }}/gatus
user=dotanarchy
autostart=true
autorestart=true
startretries=3
startsecs=1
stopcommand=bash -c "cd {{ user_home }}/gatus && podman-compose down"
stopsignal=INT
stopwaitsecs=30
stdout_logfile={{ user_home }}/.supervisor/logs/gatus.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile={{ user_home }}/.supervisor/logs/gatus.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="{{ user_home }}"
