# ansible/files/supervisord/conf.d/jellyfin.ini
[program:jellyfin]
command=/usr/bin/jellyfin -d {{ user_home }}/.config/jellyfin/data -C {{ user_home }}/.config/jellyfin/cache -c {{ user_home }}/.config/jellyfin/config -l {{ user_home }}/.config/jellyfin/log -w /usr/share/jellyfin/jellyfin-web --ffmpeg /usr/bin/ffmpeg
user=dotanarchy
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile={{ user_home }}/.supervisor/logs/jellyfin.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile={{ user_home }}/.supervisor/logs/jellyfin.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
