# ansible/files/supervisord/conf.d/karakeep.ini
[program:karakeep]
command=bash -c "cd {{ user_home }}/karakeep && podman-compose down --remove-orphans 2>/dev/null; podman-compose up -d; exec sleep infinity"
directory={{ user_home }}/karakeep
user=dotanarchy
autostart=true
autorestart=true
startretries=3
startsecs=1
stopcommand=bash -c "cd {{ user_home }}/karakeep && podman-compose down"
stopsignal=INT
stopwaitsecs=30
stdout_logfile={{ user_home }}/.supervisor/logs/karakeep.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile={{ user_home }}/.supervisor/logs/karakeep.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="{{ user_home }}"
