# ansible/files/supervisord/conf.d/n8n.ini
[program:n8n]
command=bash -c "cd {{ user_home }}/n8n && podman-compose down 2>/dev/null; podman pod rm -f n8n 2>/dev/null; podman-compose up -d"
directory={{ user_home }}/n8n
user=dotanarchy
autostart=true
autorestart=false
startretries=3
startsecs=1
stopcommand=bash -c "cd {{ user_home }}/n8n && podman-compose down"
stopsignal=INT
stopwaitsecs=30
stdout_logfile={{ user_home }}/.supervisor/logs/n8n.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile={{ user_home }}/.supervisor/logs/n8n.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="{{ user_home }}"
