# ansible/files/supervisord/conf.d/qbwc-bridge.ini
[program:qbwc-bridge]
command=bash -c "cd {{ user_home }}/qbwc-bridge && podman-compose down 2>/dev/null; podman pod rm -f qbwc-bridge 2>/dev/null; podman-compose up --build -d"
directory={{ user_home }}/qbwc-bridge
user=dotanarchy
autostart=true
autorestart=false
startretries=3
startsecs=1
stopcommand=bash -c "cd {{ user_home }}/qbwc-bridge && podman-compose down"
stopsignal=INT
stopwaitsecs=30
stdout_logfile={{ user_home }}/.supervisor/logs/qbwc-bridge.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile={{ user_home }}/.supervisor/logs/qbwc-bridge.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="{{ user_home }}"
