# ansible/files/supervisord/conf.d/tiny-stats.ini
[program:tiny-stats]
command=/usr/bin/node {{ user_home }}/tiny-stats/server.js
directory={{ user_home }}/tiny-stats
user=dotanarchy
autostart=true
autorestart=true
startretries=3
startsecs=5
environment=PORT="{{ tiny_stats_port }}"
stdout_logfile={{ user_home }}/.supervisor/logs/tiny-stats.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile={{ user_home }}/.supervisor/logs/tiny-stats.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
